# engine/launchd/com.coinbot.engine-4h.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.coinbot.engine-4h</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/riss/project/coin-bot-4h/engine/.venv/bin/python</string>
        <string>-m</string>
        <string>engine.main</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/riss/project/coin-bot-4h/engine</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/riss/project/coin-bot-4h/engine/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/riss/project/coin-bot-4h/engine/logs/launchd.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/riss/project/coin-bot-4h/engine/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
